# enolemutta.ini
[program:enolemutta]
command=/home/ubuntu/scriptapps/perkele/enolemutta.py
redirect_stderr=true
stdout_logfile=/var/log/enolemutta.out
