# kysymyksia.ini
[program:kysymyksia]
command=/home/ubuntu/scriptapps/perkele/kysymyksia.py
redirect_stderr=true
stdout_logfile=/var/log/kysymyksia.out
