# perkele.ini
[program:perkele]
command=/home/ubuntu/scriptapps/perkele/perkele.py
redirect_stderr=true
stdout_logfile=/var/log/perkele.out
